# etc/supervisord.d/slack-relay.ini
[program:slack-relay]
directory = /opt/qnib/slack-relay/
command = /opt/qnib/bin/start_slack-relay.py
stdout_logfile = /var/log/supervisor/%(program_name)s.log
redirect_stderr=true
autorestart = true
